# etc/xinetd.d/check_marathon
# default: on
# description: Check Marathon service from Haproxy to ensure layer 7 health
service check_marathon
{
	socket_type		= stream
	port            	= 60210
	protocol		= tcp
	wait			= no
	user			= root
	server			= /usr/local/bin/check_marathon.py
	server_args		= -n mesos-m11 -p 8080
	disable			= no
	per_source		= UNLIMITED
	cps			= 100 2
	flags			= REUSE
	only_from		= 0.0.0.0/0
	log_on_failure  	+= USERID
}
