# etc/xinetd.d/check_mesos
# default: on
# description: Check Mesos service from Haproxy to ensure layer 7 health
service check_mesos
{
	socket_type		= stream
	port            	= 60211
	protocol		= tcp
	wait			= no
	user			= root
	server			= /usr/local/bin/check_mesos.py
	server_args		=  -n mesos-m11 -p 5050
	disable			= no
	per_source		= UNLIMITED
	cps			= 100 2
	flags			= REUSE
	only_from		= 0.0.0.0/0
	log_on_failure  	+= USERID
}
